# data/supervisord/conf.d/autoEmailWork.ini
[program:autoEmailWork]
process_name = %(program_name)s_%(process_num)02d
command = php /www/wwwroot/pay.iredcap.cn/think queue:work --daemon --queue AutoEmailWrok
autostart = true
autorestart = true
user = www
numprocs = 1
redirect_stderr = true
stdout_logfile = /www/wwwroot/pay.iredcap.cn/data/supervisord/email_work.log